# ubuntu/bluetooth-touch-vince.conf
description     "Xiaomi redmi 5 plus bluetooth initialization"

start on starting bluetooth

expect fork
exec /bin/sh /system/halium/usr/share/bluetooth-touch/vince
